# automation/launchd/com.skippalgo.c13.ibkr-smoke.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.skippalgo.c13.ibkr-smoke</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <!-- REPLACE __REPO_PATH__ with the absolute path to your local
             skipp-algo checkout when copying this plist to
             ~/Library/LaunchAgents/. See automation/launchd/README.md. -->
        <string>__REPO_PATH__/automation/launchd/run-c13-ibkr-smoke.sh</string>
    </array>

    <!-- 08:00 local time (ET) Mon-Fri — 90 minutes before US-equity open
         (09:30 ET). TWS Paper Gateway must already be running on 127.0.0.1:7497.
         On EXIT != 0 the driver writes cache/live/smoke_HALT; the execution
         runner (run_ibkr_open_execution.py) refuses to submit orders if that
         sentinel file is present or if no smoke JSONL for today exists. -->
    <key>StartCalendarInterval</key>
    <array>
        <dict><key>Weekday</key><integer>1</integer><key>Hour</key><integer>8</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Weekday</key><integer>2</integer><key>Hour</key><integer>8</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Weekday</key><integer>3</integer><key>Hour</key><integer>8</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Weekday</key><integer>4</integer><key>Hour</key><integer>8</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Weekday</key><integer>5</integer><key>Hour</key><integer>8</integer><key>Minute</key><integer>0</integer></dict>
    </array>

    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/tmp/skippalgo-c13-ibkr-smoke.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/skippalgo-c13-ibkr-smoke.err</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
