# automation/launchd/com.skippalgo.c13.phase-a-export.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.skippalgo.c13.phase-a-export</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <!-- REPLACE __REPO_PATH__ with the absolute path to your local
             skipp-algo checkout when copying this plist to
             ~/Library/LaunchAgents/. See automation/launchd/README.md. -->
        <string>__REPO_PATH__/automation/launchd/run-c13-phase-a-export.sh</string>
    </array>

    <!-- 09:18 local time (ET) Mon-Fri — 10 minutes BEFORE the phase-a
         runner so fresh trade-cards land before build_phase_a_inputs.py
         picks them up. -->
    <key>StartCalendarInterval</key>
    <array>
        <dict><key>Weekday</key><integer>1</integer><key>Hour</key><integer>9</integer><key>Minute</key><integer>18</integer></dict>
        <dict><key>Weekday</key><integer>2</integer><key>Hour</key><integer>9</integer><key>Minute</key><integer>18</integer></dict>
        <dict><key>Weekday</key><integer>3</integer><key>Hour</key><integer>9</integer><key>Minute</key><integer>18</integer></dict>
        <dict><key>Weekday</key><integer>4</integer><key>Hour</key><integer>9</integer><key>Minute</key><integer>18</integer></dict>
        <dict><key>Weekday</key><integer>5</integer><key>Hour</key><integer>9</integer><key>Minute</key><integer>18</integer></dict>
    </array>

    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>__HOME__/Library/Logs/skippalgo/skippalgo-c13-phase-a-export.log</string>
    <key>StandardErrorPath</key>
    <string>__HOME__/Library/Logs/skippalgo/skippalgo-c13-phase-a-export.err</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
