# automation/launchd/com.skippalgo.c13.phase-a-prep.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.skippalgo.c13.phase-a-prep</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <!-- REPLACE __REPO_PATH__ with the absolute path to your local
             skipp-algo checkout when copying this plist to
             ~/Library/LaunchAgents/. See automation/launchd/README.md. -->
        <string>__REPO_PATH__/automation/launchd/run-c13-phase-a-prep.sh</string>
    </array>

    <!-- 09:15 local time (ET on the Phase-A workstation) Mon-Fri, ~15
         minutes before the cash open. Writes setups_<DATE>.jsonl and
         gate_status_<DATE>.json into cache/phase_a/ so the session
         agent can launch with up-to-date inputs at 09:25.
         StartCalendarInterval is local-clock so DST is handled by the OS. -->
    <key>StartCalendarInterval</key>
    <array>
        <dict><key>Weekday</key><integer>1</integer><key>Hour</key><integer>9</integer><key>Minute</key><integer>15</integer></dict>
        <dict><key>Weekday</key><integer>2</integer><key>Hour</key><integer>9</integer><key>Minute</key><integer>15</integer></dict>
        <dict><key>Weekday</key><integer>3</integer><key>Hour</key><integer>9</integer><key>Minute</key><integer>15</integer></dict>
        <dict><key>Weekday</key><integer>4</integer><key>Hour</key><integer>9</integer><key>Minute</key><integer>15</integer></dict>
        <dict><key>Weekday</key><integer>5</integer><key>Hour</key><integer>9</integer><key>Minute</key><integer>15</integer></dict>
    </array>

    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/tmp/skippalgo-c13-phase-a-prep.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/skippalgo-c13-phase-a-prep.err</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
